# recipes-ve/dbus-modbustcp/files/service/log/run
#!/bin/sh
exec 2>&1
exec multilog s99999 n8 /log/dbus-modbustcp
