# meta-venus/recipes-extended/haveged/haveged.rules
ACTION=="add", KERNEL=="random", RUN+="/usr/sbin/haveged -F -w 1024 --once"
