# meta-venus/recipes-core/udev/files/bluetooth.rules
ACTION=="add", SUBSYSTEM=="bluetooth", ENV{DEVTYPE}=="host", RUN+="bt-config %k"
ACTION=="remove", SUBSYSTEM=="bluetooth", ENV{DEVTYPE}=="host", RUN+="bt-remove %k"
