# meta-venus/recipes-core/udev/files/rtl8192cu.rules
# workaround for replugging a wifu usb adapter
# with unload the module, it won't work when plugging it in again
ACTION=="remove", ENV{ID_BUS}=="usb", ENV{ID_USB_DRIVER}=="rtl8192cu", RUN+="/sbin/rmmod rtl8192cu"
