# meta-venus/recipes-core/initscripts/files/usb.rules
ACTION=="add", KERNEL=="sd?[0-9]", SUBSYSTEM=="block", RUN+="/etc/init.d/update-data.sh -r"
